# python/logging/logging.crontab
0 * * * * /usr/local/stork/bin/storklog.py &> /dev/null
